# pi/motion-display-controller/service/motion-display.service
[Unit]
Description=Motion Display Controller
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/motion-display-controller
ExecStart=/usr/bin/python3 /opt/motion-display-controller/motion_display.py
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

# Optional: safer shutdown handling
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target